# etc/supervisor.conf
[program:standing-desk]
command=python /home/pi/standing-desk/src/main.py
directory=/home/pi/standing-desk/src/
user=pi
autostart=true
autorestart=true
stopsignal=INT
priority=999
